# debian-12-bookworm-std-64/std_setup
#! /bin/sh
### BEGIN INIT INFO
# Provides:          std_setup
# Required-Start:    $local_fs
# Required-Stop:
# X-Start-Before:
# Default-Start:     2
# Default-Stop:
# Short-Description: Generic 64 bit Bookworm Standard Installer
# Description: Generic 64 bit Bookworm Standard instance installer invoked during first boot and then removes itself
### END INIT INFO

# Author: Ap.Muthu <apmuthu@usa.net>

set -e

# DOMAIN=`grep '^search' /etc/resolv.conf|awk '{ print $2; }'`
HNAME=`head -n 1 /etc/hostname|awk '{ print $1; }'`

if [ "X${HNAME}" = "Xlocalhost" ] ; then
	exit 0;
fi

# Remove delay in login
systemctl mask systemd-logind

# now certs are autogenerated
# rm -rf /root/.ssh
# ssh-keygen -q -t rsa -b 4096 -f ~/.ssh/id_rsa -P "" <<<y >/dev/null 2>&1

if [ -x /sbin/insserv ] ; then
    /sbin/insserv -r std_setup
    rm -f /etc/init.d/std_setup
else
    rm -f /etc/init.d/std_setup
    update-rc.d -f std_setup remove
fi

